# systemd/hourly-time-bot.service
[Unit]
Description=Hourly Time Bot for MAX Platform
Documentation=https://github.com/ramtuk/hourly-time-bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=voiceai
Group=voiceai
WorkingDirectory=/opt/hourly-time-bot
Environment=PATH=/opt/hourly-time-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/opt/hourly-time-bot
EnvironmentFile=/opt/hourly-time-bot/.env
ExecStart=/opt/hourly-time-bot/venv/bin/python /opt/hourly-time-bot/run.py

# Логирование в /var/log
StandardOutput=file:/var/log/hourly-time-bot.log
StandardError=file:/var/log/hourly-time-bot.error.log

# Ротация логов (systemd будет управлять ротацией)
LogsDirectory=hourly-time-bot
LogsDirectoryMode=0755

# Process management
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
